# scripts/paperman-server.service
[Unit]
Description=Paperman Search Server
Documentation=https://github.com/sjg20/paperman
After=network.target

[Service]
Type=simple
User=paperman
Group=paperman
WorkingDirectory=/opt/paperman
ExecStart=/opt/paperman/paperman-server /srv/papers
Restart=on-failure
RestartSec=5s

# Security settings
NoNewPrivileges=true
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=paperman-server

[Install]
WantedBy=multi-user.target
